# bash_ex/task_4/write_time.service
[Unit]
Description=write time every second

[Service]
Type=oneshot
StartLimitBurst=0
ExecStart=/home/ivan/Desktop/task_4/write_time.sh
