# bash_ex/task_4/write_time.timer
[Unit]

[Timer]
AccuracySec=1us
OnUnitActiveSec=1s
OnUnitInactiveSec=1s
OnBootSec=1s
Unit=write_time.service

[Install]
WantedBy=multi-user.target
