# bash_ex/task_4/write_time.timer
[Unit]

[Timer]
AccuracySec=1us
OnCalendar=*:*:0/1
Unit=write_time.service

[Install]
WantedBy=timers.target
