# bash_ex/task_4/write_time.timer
[Unit]

[Timer]
OnBootSec=1
OnCalendar=*:*:0/1
Unit=write_time.service

[Install]
WantedBy=multi-user.target